# systemd/swaybg-swapper.timer
[Unit]
Description=Executes swaybg-swapper.service
[Timer]
OnUnitInactiveSec=60m
[Install]
WantedBy=timers.target
